# etc/qualys/cloud-agent/.systemd/qualys-cloud-agent-uninstall.service
[Unit]
Description=Qualys cloud agent uninstall

[Service]
Type=simple
ExecStart=/usr/bin/rpm -ev qualys-cloud-agent
